# cookbooks/install_vault/templates/vault.service.erb
[Unit]
Description=Vault Secret Encryption Service
Requires=network-online.target
After=network.target
Documentation=https://www.vaultproject.io/docs/

[Service]
Restart=on-failure
Environment=GOMAXPROCS=2 "VAULT_ADDR='http://127.0.0.1:8200'"
ExecStartPre=/bin/export $VAULT_ADDR
ExecStart=/opt/vault/vault server -config=/opt/vault/config.hcl
ExecReload=/bin/kill -s HUP $MAINPID
KillSignal=SIGINT
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
